# systemd/isolate.service.in
[Unit]
Description=A trivial daemon to keep Isolate's control group hierarchy

[Service]
Type=notify
ExecStart=@SBINDIR@/isolate-cg-keeper
Slice=isolate.slice
Delegate=true

[Install]
WantedBy=multi-user.target
