# systemd/isolate.slice
[Unit]
Description=Slice for Isolate's sandboxes
